# etc/rc.local
#!/bin/sh

set -e

systemctl restart fix-sound.service

exit 0
